# fb-ads-scraper/com.user.elixir-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.elixir-dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>http.server</string>
        <string>8080</string>
        <string>--bind</string>
        <string>127.0.0.1</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/vaceslavlusnikov/Desktop/дашборд</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/elixir_dashboard.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/elixir_dashboard.err</string>
</dict>
</plist>
