# fb-ads-scraper/com.user.fbscraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.fbscraper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vaceslavlusnikov/Desktop/дашборд/fb-ads-scraper/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vaceslavlusnikov/Desktop/дашборд/fb-ads-scraper</string>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/fb_ads_scraper.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/fb_ads_scraper.err</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
